# dashboard_server/supervisor.conf
[program:hbp_dashboard]
user=hbp_dashboard
directory=/opt/hbp_dashboard/hotel_book_prob
command=/opt/hbp_dashboard/.virtualenvs/hbp_dashboard/bin/python /opt/hbp_dashboard/.virtualenvs/hbp_dashboard/bin/gunicorn --name=hbp_dashboard --bind=0.0.0.0:5000 --workers=2 dashboard_server.app:app --timeout 90 --log-file /var/log/hbp_dashboard/hbp_dashboard_gunicorn.log
redirect_stderr=true
stdout_logfile=/var/log/supervisor/hbp_dashboard.log
stderr_logfile=/var/log/supervisor/hbp_dashboard.err
stdout_logfile_maxbytes=5MB
stderr_logfile_maxbytes=5MB
stdout_logfile_backups=1
stderr_logfile_backups=1
stopasgroup=true
autorestart=true
exitcodes=0,2